# staged/services/context-intelligence-unified/context-intelligence-unified.service
[Unit]
Description=Syn_OS Unified Context Intelligence Service
Documentation=https://github.com/syn-os/syn-os
After=network.target consciousness-unified.service
Wants=network-online.target

[Service]
Type=simple
User=synos
Group=synos
WorkingDirectory=/opt/synos/services/context-intelligence-unified
ExecStart=/usr/bin/python3 unified_context_intelligence_service.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=context-intelligence-unified

# Environment variables
Environment=PYTHONPATH=/opt/synos/services/context-intelligence-unified
Environment=PYTHONUNBUFFERED=1

# Resource limits
MemoryMax=4G
CPUQuota=300%

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/synos/services/context-intelligence-unified/logs
ReadWritePaths=/opt/synos/services/context-intelligence-unified/data
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target