# operations/admin/syn-enterprise-dashboard.service
[Unit]
Description=Syn_OS Enterprise Dashboard
After=network.target syn-consciousness.service
Wants=network.target
Requires=syn-consciousness.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/syn_os
Environment=PYTHONPATH=/opt/syn_os
ExecStart=/usr/bin/python3 -m src.enterprise.dashboard
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/syn_os/data /var/log
PrivateTmp=true

[Install]
WantedBy=multi-user.target
