# services/consciousness-dashboard/consciousness-dashboard.service
[Unit]
Description=SynapticOS Consciousness Dashboard
After=network.target consciousness-ai-bridge.service
Wants=consciousness-ai-bridge.service

[Service]
Type=simple
User=synos
Group=synos
WorkingDirectory=/opt/synapticos/services/consciousness-dashboard
ExecStart=/usr/bin/python3 consciousness_dashboard.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/synapticos
Environment=DASHBOARD_PORT=8000
Environment=CONSCIOUSNESS_BRIDGE_URL=http://localhost:8082
Environment=EDUCATION_SERVICE_URL=http://localhost:8084

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=consciousness-dashboard

[Install]
WantedBy=multi-user.target
