# workspace/integration/parrotos/overlay/systemd-services/synapticos-news-intelligence.service
[Unit]
Description=SynapticOS News Intelligence Platform
After=network.target synapticos-consciousness-bridge.service
Requires=synapticos-consciousness-bridge.service

[Service]
Type=simple
User=synapticos
Group=synapticos
WorkingDirectory=/opt/synapticos/services/news-intelligence
ExecStart=/usr/bin/python3 news_intelligence.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/synapticos/services/news-intelligence

[Install]
WantedBy=multi-user.target
