# workspace/integration/parrotos/overlay/systemd-services/synapticos-consciousness-bridge.service
[Unit]
Description=SynapticOS Consciousness AI Bridge
After=network.target
Wants=postgresql.service redis.service

[Service]
Type=simple
User=synapticos
Group=synapticos
WorkingDirectory=/opt/synapticos/services/consciousness-ai-bridge
ExecStart=/usr/bin/python3 consciousness_bridge.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/synapticos/services/consciousness-ai-bridge

[Install]
WantedBy=multi-user.target
