# linux-distribution/SynOS-Linux-Builder/config/includes.chroot/etc/xdg/autostart/synos-firstboot.desktop
[Desktop Entry]
Type=Application
Name=SynOS First-Boot Wizard
Comment=Initial configuration wizard for SynOS v1.0
Exec=/usr/local/bin/synos-firstboot-wizard
Terminal=true
Categories=System;Settings;
X-MATE-Autostart-enabled=true
Hidden=false
NoDisplay=false
OnlyShowIn=MATE;
X-MATE-Autostart-Phase=Application
