# linux-distribution/SynOS-Linux-Builder/config/includes.chroot/etc/xdg/autostart/synos-ai-chat-panel.desktop
[Desktop Entry]
Type=Application
Name=SynOS AI Chat Panel
Comment=Jarvis AI Assistant Chat Interface
Exec=/usr/local/bin/synos-ai-chat-panel
Icon=applications-chat
Terminal=false
Categories=Utility;System;
StartupNotify=false
X-MATE-Autostart-enabled=true
Hidden=false
