# archive-consolidated/05-legacy-systems/legacy-v0.x/services/consciousness-unified/consciousness-unified.service
[Unit]
Description=Syn_OS Unified Consciousness Service
Documentation=https://github.com/syn-os/syn-os
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=synos
Group=synos
WorkingDirectory=/opt/synos/services/consciousness-unified
ExecStart=/usr/bin/python3 unified_consciousness_service.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=consciousness-unified

# Environment variables
Environment=PYTHONPATH=/opt/synos/services/consciousness-unified
Environment=PYTHONUNBUFFERED=1

# Resource limits
MemoryMax=2G
CPUQuota=200%

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/synos/services/consciousness-unified/logs
ReadWritePaths=/tmp

[Install]
WantedBy=multi-user.target