# parrotos-integration/overlay/services/educational-platform/education-gui.service
[Unit]
Description=SynapticOS Educational Platform GUI
After=network.target educational-platform.service consciousness-ai-bridge.service
Wants=educational-platform.service consciousness-ai-bridge.service

[Service]
Type=simple
User=synos
Group=synos
WorkingDirectory=/opt/synapticos/services/educational-platform
ExecStart=/usr/bin/python3 education_gui.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/synapticos
Environment=GUI_PORT=8001
Environment=EDUCATION_SERVICE_URL=http://localhost:8084
Environment=CONSCIOUSNESS_BRIDGE_URL=http://localhost:8082

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=education-gui

[Install]
WantedBy=multi-user.target
