# parrotos-integration/overlay/systemd-services/synapticos-ctf-platform.service
[Unit]
Description=SynapticOS CTF Platform
After=network.target synapticos-consciousness-bridge.service
Requires=synapticos-consciousness-bridge.service

[Service]
Type=simple
User=synapticos
Group=synapticos
WorkingDirectory=/opt/synapticos/services/ctf-platform
ExecStart=/usr/bin/python3 ctf_platform.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/synapticos/services/ctf-platform

[Install]
WantedBy=multi-user.target
