# parrotos-integration/overlay/systemd-services/synapticos-context-engine.service
[Unit]
Description=SynapticOS Advanced Context Engine
After=network.target synapticos-consciousness-bridge.service
Requires=synapticos-consciousness-bridge.service

[Service]
Type=simple
User=synapticos
Group=synapticos
WorkingDirectory=/opt/synapticos/services/context-engine
ExecStart=/usr/bin/python3 context_engine.py
Restart=always
RestartSec=10
Environment=PYTHONPATH=/opt/synapticos/services/context-engine

[Install]
WantedBy=multi-user.target
