# archives/2025-10/legacy-builders/SynOS-Linux-Builder/systemd-services/synos-behavior-monitor.service
[Unit]
Description=SynOS Behavioral Monitor Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/synos-apps/behavior
ExecStart=/usr/bin/python3 /opt/synos-apps/behavior/behavior_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
