# archives/2025-10/legacy-builders/SynOS-Linux-Builder/systemd-services/synos-consciousness-integration.service
[Unit]
Description=SynOS Consciousness Integration Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/synos-apps/consciousness
ExecStart=/usr/bin/python3 /opt/synos-apps/consciousness/consciousness_bridge.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
